# muxer/runtime-package/systemd/ike-nat-bridge.service
[Unit]
Description=MUXER2 NFQUEUE IKE Rewrite Engine (Marker + NATD DPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /etc/muxer/src/ike_nat_bridge.py --log-level INFO
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
