# muxer/runtime-package/systemd/muxer-trace.service
[Unit]
Description=MUXER2 Packet Trace (pcap ring + counter snapshots)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /etc/muxer/src/mux_trace.py --log-level INFO
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
