# scripts/customers/systemd/rpdb-jump-host-route-sync.service
[Unit]
Description=RPDB Jump Host Customer Route Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
Group=ec2-user
WorkingDirectory=/home/ec2-user/rpdb
Environment=RPDB_ENVIRONMENT=rpdb-empty-live
ExecStart=/usr/bin/python3 scripts/customers/sync_jump_host_customer_routes.py --environment ${RPDB_ENVIRONMENT} --follow
Restart=always
RestartSec=15

[Install]
WantedBy=multi-user.target
