# scripts/customers/systemd/rpdb-dynamic-peer-ip-watcher.service
[Unit]
Description=RPDB Dynamic Peer IP Watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/muxingRPDB-Platform-Framework
Environment=RPDB_ENVIRONMENT=rpdb-empty-live
ExecStart=/usr/bin/python3 scripts/customers/run_dynamic_peer_ip_watcher.py --environment ${RPDB_ENVIRONMENT} --follow --approve
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
